# systemd.service
[Unit]
Description=Forge Service

[Service]
Type=simple
User=forge
Password=forge
ExecStart=/usr/lib/forge/forged